# files/node/usr/service/proxy/run
#!/bin/bash

exec 2>&1
exec /usr/local/bin/kube-proxy  \
	--cluster-cidr=10.200.0.0/16 \
	--kubeconfig=/var/lib/kubernetes/kube-proxy.kubeconfig \
	--proxy-mode=iptables \
	--v=2

