# files/node/usr/service/dockerd/run
#!/bin/bash

exec 2>&1
exec /usr/bin/dockerd \
        --iptables=false \
        --ip-masq=false \
        --bridge="cbr0"

