# files/node/etc/network/if-up.d/route
#!/bin/sh -e

/sbin/ip route add {{ vip }}/32 dev docker0

