# files/node/usr/service/kubelet/run
#!/bin/bash

exec 2>&1
exec /usr/local/bin/kubelet  \
	--kubeconfig=/var/lib/kubelet/kubeconfig \
	--v 2 \
	--allow-privileged=true \

